# systemd/.config/systemd/user/parucache.timer
[Unit]
Description=Run paru-clean weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
