# systemd/.config/systemd/user/waybar-updates.service
[Unit]
Description=Populate update cache for Waybar

[Service]
Type=oneshot
ExecStart=%h/.scripts/update-timer.sh

